# deploy/systemd/blocky-brain.service
[Unit]
Description=Blocky signal brain
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=blocky
Group=blocky
WorkingDirectory=/opt/blocky-polymarket
EnvironmentFile=/opt/blocky-polymarket/.env
ExecStart=/opt/blocky-polymarket/.venv/bin/python -u brain/main.py
Restart=always
RestartSec=10
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/blocky-polymarket/data

[Install]
WantedBy=multi-user.target
