# deploy/systemd/blocky-executor.service
[Unit]
Description=Blocky trade executor
After=network-online.target blocky-brain.service
Wants=network-online.target

[Service]
Type=simple
User=blocky
Group=blocky
WorkingDirectory=/opt/blocky-polymarket
EnvironmentFile=/opt/blocky-polymarket/.env
ExecStart=/opt/blocky-polymarket/.venv/bin/python -m pyapp.executor
Restart=always
RestartSec=10
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/blocky-polymarket/data

[Install]
WantedBy=multi-user.target
